# backend/budgetal.service.example
[Unit]
Description=Budgetal
ConditionPathExists=/home/deploy/budgetal-beta/budgetal
After=network.target

[Service]
Environment=GO_ENV=production
Environment=DATABASE_URL=postgres:///budgetal_production?host=/var/run/postgresql&fallback_application_name=Budgetal-Go
Environment=CORS=https://example.com

Environment=SMTP_USER=sendgriduser
Environment=SMTP_PASSWORD=keep it secret, keep it safe
Environment=SMTP_HOST=smtp.sendgrid.net
Environment=SMTP_PORT=587
Environment=ERROR_NOTIFICATION_EMAILS=email@example.com,email2@example.com

Environment=AWS_S3_ACCESS_KEY_ID=keyid
Environment=AWS_S3_SECRET_ACCESS_KEY=se/cr/et
Environment=AWS_S3_TOKEN=
Environment=AWS_S3_REGION=us-east-1
Environment=AWS_S3_BUCKET=my-bucket

Type=simple
User=deploy
Group=deploy
LimitNOFILE=1024

Restart=on-failure
RestartSec=10
startLimitIntervalSec=60

WorkingDirectory=/home/deploy/budgetal-beta
ExecStart=/home/deploy/budgetal-beta/budgetal

# make sure log directory exists and owned by syslog
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/log/budgetal-beta
ExecStartPre=/bin/chown syslog:adm /var/log/budgetal-beta
ExecStartPre=/bin/chmod 755 /var/log/budgetal-beta
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=budgetal

[Install]
WantedBy=multi-user.target